# qr3_70.conf
[program:qr3_70]
command=python3 /home/curator3/qr3_70/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_70/test.err.log
stdout_logfile=/home/curator3/qr3_70/test.err.log
